# repo-man.service
[Unit]
Description=Penguin repo manager

[Service]
ExecStart=/usr/bin/repo-man
Type=simple
Restart=always

[Install]
WantedBy=default.target
